# admin/etc/dumpSlurmDB.cron
# nightly dump of SlurmDB, place in /etc/cron.d
05 02 * * * root /home/software/tools/admin/bin/dump_SlurmDB.sh
