# admin/etc/restart_ccdc-licence.cron
#!/bin/bash

#place this script as /etc/cron.hourly/restart_ccdc-licence
/bin/systemctl restart ccdc-licence